# deploy/zhihuishu-login-cleanup.service
[Unit]
Description=Canvas Dashboard browser login session cleanup
After=docker.service

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/canvas-dashboard/current
ExecStart=/home/ubuntu/canvas-dashboard/.venv/bin/python zhihuishu_login_sessions.py --cleanup-expired
ExecStart=/home/ubuntu/canvas-dashboard/.venv/bin/python tongji_login_sessions.py --cleanup-expired
